# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:dummy-llm]
command=python /app/services/dummy-llm/app.py
environment=PORT="8001"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=5

[program:agentic-ai]
command=python /app/services/agentic-ai/app.py
environment=PORT="8002",INSURANCE_AI_AGENT_URL="http://localhost:8003",DISPATCH_AI_AGENT_URL="http://localhost:8004"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=5

[program:insurance-agent]
command=python /app/services/insurance-ai-agent/app.py
environment=PORT="8003"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=5

[program:dispatch-agent]
command=python /app/services/dispatch-ai-agent/app.py
environment=PORT="8004"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=5

[program:telegram-bot]
command=python /app/services/telegram-bot/app.py
environment=LLM_SERVICE_URL="http://localhost:8001",AGENTIC_AI_URL="http://localhost:8002"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=5
startsecs=5
